# crontab
*/2 * * * * root /usr/sbin/arp -vn > /data/arp/$(hostname)_$(date +\%H-\%M-\%S) 2>&1
*/2 * * * * root /sbin/route -n > /data/routes/$(hostname)_$(date +\%H-\%M-\%S) 2>&1
#required line